# data/systemd/pygcn.service
[Unit]
Description=PyGCN service

[Service]
Nice=18
ExecStart=pygcn-listen

[Install]
WantedBy=default.target

